# systemd/account-intel-ui-internal.service
[Unit]
Description=Jazzware Account Intel — Internal UI (Streamlit, port 8502)
After=network-online.target account-intel-api.service
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/jazzware-account-intel
EnvironmentFile=/opt/jazzware-account-intel/.env
ExecStart=/opt/jazzware-account-intel/.venv/bin/streamlit run src/account_intel/ui/internal_app.py --server.address 0.0.0.0 --server.port 8502 --server.headless true --browser.gatherUsageStats false
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
